# deploy/prod/backup/systemd/para-backup.service
[Unit]
Description=para backup (pg_dump + assets snapshot)
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/dockers/para/app

Environment=BACKUP_ROOT=/root/dockers/para/backups
Environment=COMPOSE_PROJECT_NAME=para
Environment=COMPOSE_FILE=/root/dockers/para/app/deploy/prod/docker-compose.yml
Environment=ENV_FILE=/root/dockers/para/.env
Environment=ASSETS_DIR=/root/dockers/para/data/server/.data

ExecStart=/usr/bin/env bash -c 'bash deploy/prod/backup/backup_pg_dump.sh && bash deploy/prod/backup/backup_assets.sh'

Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
